# systemd/netmuxd.service
[Unit]
Description=netmuxd for AltServer
After=network-online.target

[Service]
ExecStart=/opt/altserver/netmuxd

[Install]
WantedBy=default.target
